# package/arduino/luci-app-arduino-webpanel/files/etc/uci-defaults/arduino-webpanel
#!/bin/sh

mv /www/arduino.index.html /www/index.html
